# deploy/roles/node-server/templates/server.service.j2
# /etc/systemd/system/myapp.service
[Service]
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /usr/bin/timberbuddy/src/server/server.js