# deploy/roles/node-server/templates/TimberBuddyInterfaceService.service.j2
[Unit]
Description=Timber Buddy Control Interface
After=systemd-udevd.target

[Service]
Type=simple
ExecStart=/usr/bin/npx tsx server/pi_socket_server.ts
WorkingDirectory=/home/{{uname}}/timberbuddy/
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target