# deploy/roles/node-server/templates/TimberBuddyInterfaceService.service.j2
[Unit]
Description=Timber Buddy Control Interface
After=ssh.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/jsullivan2/timberbuddy/server/index.js
WorkingDirectory=/home/jsullivan2/timberbuddy/
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target